# deploy/etc/cliche-celery-beat.conf
description "Cliche celery beat"

start on (local-filesystems and runlevel [2345])
stop on runlevel [016]

respawn

exec /home/cliche/bin/celery -A cliche -c /home/cliche/etc/prod.cfg.yml