# deploy/etc/cliche-celery-beat.conf
description "Cliche celery beat"

start on (local-filesystems and runlevel [2345])
stop on runlevel [016]

exec /home/cliche/bin/celery beat -A cliche --config /home/cliche/etc/prod.cfg.yml
